# packaging/mediacontroller.service
[Unit]
Description=Media controller

[Service]
ExecStart=/usr/bin/mediacontroller
Type=simple

[Install]
Also=mediacontroller.socket
