# packaging/mediacontroller-ipc.socket
[Unit]
Description=MediaController Service socket

[Socket]
ListenStream=/run/.mediacontroller-comm.sock

[Install]
WantedBy=sockets.target
