# packaging/mediacontroller.socket
[Unit]
Description=MediaController Service socket

[Socket]
SocketUser=multimedia_fw
SocketGroup=multimedia_fw
ListenStream=/var/run/media-controller/media_sa_controller
Service=mediacontroller.service

[Install]
WantedBy=sockets.target
